# rules/80-local.rules
SUBSYSTEM=="block",KERNEL == "sda1" ,ATTRS{vendor}=="Kingston",ACTION=="add", RUN+="/usr/local/bin/trigger.sh"